# root/etc/services.d/030_nginx/run
#!/bin/bash

exec /opt/nginx/sbin/nginx -c /etc/nginx.conf
