# root/etc/services.d/030_nginx/run
#!/bin/bash

exec /opt/nginx/sbin/nginx -c /opt/nginx/conf/nginx.conf
